# contrib/systemd/phantasm.service
[Unit]
Description=Phantasm local WebUI
After=network.target

[Service]
Type=simple
User=phantasm
Group=phantasm
WorkingDirectory=/opt/phantasm
Environment=PYTHONPATH=/opt/phantasm/src
Environment=PHANTASM_FIELD_MODE=1
Environment=PHANTASM_AUDIT=0
Environment=PHANTASM_DEBUG=0
Environment=PHANTASM_HOST=127.0.0.1
Environment=PHANTASM_PORT=8000
Environment=PHANTASM_STATE_DIR=/var/lib/phantasm
ExecStart=/usr/bin/python3 -m phantasm.web_server
Restart=on-failure
RestartSec=2
UMask=0077

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/phantasm /var/lib/phantasm /run/phantasm
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallFilter=@system-service
MemoryDenyWriteExecute=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
