# contrib/systemd/phasmid.service
[Unit]
Description=Phasmid local WebUI
After=network.target

[Service]
Type=simple
User=phasmid
Group=phasmid
WorkingDirectory=/opt/phasmid
Environment=PYTHONPATH=/opt/phasmid/src
Environment=PHASMID_FIELD_MODE=1
Environment=PHASMID_AUDIT=0
Environment=PHASMID_DEBUG=0
Environment=PHASMID_HOST=127.0.0.1
Environment=PHASMID_PORT=8000
Environment=PHASMID_STATE_DIR=/var/lib/phasmid
ExecStart=/usr/bin/python3 -m phasmid.web_server
Restart=on-failure
RestartSec=2
UMask=0077

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/phasmid /var/lib/phasmid /run/phasmid
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallFilter=@system-service
MemoryDenyWriteExecute=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
